# modules/starship/.config/profile.d/starship.sh
#
# Initialize starship
#

# XDG specification
export STARSHIP_CONFIG="${XDG_CONFIG_HOME}/starship/starship.toml"

export STARSHIP_CACHE="${XDG_CACHE_HOME}/startship"